# config/autostart/terminal-tabs.desktop
[Desktop Entry]
Type=Application
Name=Terminal with Tabs
Exec=bash /home/j/gbin/launch_terminal_tabs.sh
Hidden=false
NoDisplay=false
X-XFCE-Autostart=true

